# Binaries/Test_Extraction_Firmware/_miwifi_r1cm_firmware_36bfd_2.22.21.bin.extracted/squashfs-root/etc/hotplug.d/iface/99-captive_portal
#!/bin/sh

if [ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ]; then
    [ -f /usr/sbin/captive_portal.sh ] || return
    logger -t "captive_portal" "reload captive portal status"
    # reload captive portal status
    /usr/sbin/captive_portal.sh reload
fi
